# rootfs/etc/update-motd.d/00-header
#!/bin/bash
cat << 'EOF'
  ______              _____ _
 |___  /             / ____| |
    / / ___ _ __ ___| |    | | __ ___      __
   / / / _ \ '__/ _ \ |    | |/ _` \ \ /\ / /
  / /_|  __/ | | (_) | |___| | (_| |\ V  V /
 /_____\___|_|  \___/ \_____|_|\__,_| \_/\_/

  ZeroClaw App Platform Container
EOF
