# rootfs/etc/update-motd.d/80-zeroclaw
#!/bin/bash
source /etc/update-motd.d/lib.sh

echo -e "  пїЅ ${BOLD}ZeroClaw${NC}"
echo -e "  ${DIM}в”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђ${NC}"

STATUS=$(zeroclaw status 2>/dev/null)
if [ $? -eq 0 ]; then
  printf "  ${CYAN}%-12s${NC} ${GREEN}вњ“ running${NC}\n" "Status:"
else
  printf "  ${CYAN}%-12s${NC} ${RED}вњ— not running${NC}\n" "Status:"
fi

echo ""
